# systemd/tracker.service
[Unit]
Description=STBet Cricket-SuperOver Continuous Tracker
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/STBet-Tracker
ExecStart=/home/ubuntu/STBet-Tracker/.venv/bin/python tracker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stbet-tracker

[Install]
WantedBy=multi-user.target
